# deployment/uberspace/redis.ini
[program:redis]
command=redis-server %(ENV_HOME)s/.redis/conf
autostart=yes
autorestart=yes
