# deployment/uberspace/rundfunk-frontend.ini
[program:rundfunk-frontend]
directory=%(ENV_HOME)s/rundfunk-frontend
command=node production-server
environment=PORT=4200,NODE_ENV=production
autostart=yes
autorestart=yes
